# skeleton/selenium-xvfb/DEBIAN/postinst
#! /bin/sh

set -e

if [ -x "/etc/init.d/selenium-xvfb" ]; then
        update-rc.d selenium-xvfb defaults 10 >/dev/null || exit $?
fi